# config/autostart/dunst.desktop
[Desktop Entry]
Name=dunst
Exec=/usr/bin/dunst
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
