# config/autostart/eww.desktop
[Desktop Entry]
Name=eww
Exec=sh -c '/usr/bin/eww daemon; /usr/bin/eww open bar'
Terminal=false
Type=Application
