# config/autostart/gammastep.desktop
[Desktop Entry]
Name=gammastep
Exec=/usr/bin/gammastep
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
