# config/autostart/ironbar.desktop
[Desktop Entry]
Name=ironbar
Exec=/usr/bin/ironbar
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
