# config/autostart/scc.desktop
[Desktop Entry]
Name=scc
Exec=/usr/bin/scc-daemon start
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
