# config/autostart/start.desktop
[Desktop Entry]
Name=start
Exec=/home/sanpi/.config/sway/start
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
