# config/autostart/udiskie.desktop
[Desktop Entry]
Name=udiskie
Exec=/usr/bin/udiskie --tray
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=true
